# wpeframework-services.path
#This file keeps track of path - whenever it's changed it will trigger thunderstartupservices.
#In thunderstartup service files we use ConditionPathExists to check whether path exixts or not.

[Unit]
Description=Path Changed
OnFailure=path-fail-notifier@%n.service

[Path]
PathChanged=/tmp/wpeframeworkstarted
Unit=wpeframework-services.target

[Install]
WantedBy=multi-user.target