# systemd/system/wpeframework-deviceinfo.service
# WPEFramework DeviceInfo Initialiser
# Activates the DeviceInfo plugin after framework initialization
[Unit]
Description=WPEFramework DeviceInfo Initialiser
Requires=wpeframework.service iarmbusd.service dsmgr.service
After=wpeframework.service iarmbusd.service dsmgr.service
ConditionPathExists=/tmp/wpeframeworkstarted
[Service]
Type=oneshot
RemainAfterExit=yes
TimeoutStartSec=30
TimeoutStopSec=10
PrivateTmp=yes
ExecStart=/usr/bin/PluginActivator DeviceInfo

[Install]
WantedBy=multi-user.target
