# systemd/system/wpeframework-ocdm.service
[Unit]
Description=WPEFramework OCDM Plugin
Requires=wpeframework.service
After=wpeframework.service wpeframework-provisioning-ready.target
Requisite=wpeframework-provisioning-ready.target
ConditionPathExists=/tmp/wpeframeworkstarted

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/PluginActivator org.rdk.OCDM

[Install]
WantedBy=wpeframework-provisioning-ready.target
