# systemd/system/wpeframework-system.service
[Unit]
Description=WPEFramework System Initialiser
Requires=wpeframework-powermanager.service sysmgr.service mfrmgr.service dsmgr.service wpeframework-avinput.service
After=wpeframework-powermanager.service sysmgr.service mfrmgr.service dsmgr.service wpeframework-avinput.service
ConditionPathExists=/tmp/wpeframeworkstarted
[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/PluginActivator org.rdk.System
