# systemd/system/wpeframework-monitor.service
# WPEFramework Monitor Initialiser
# Activates the Monitor plugin after framework initialization
[Unit]
Description=WPEFramework Monitor Initialiser
Requires=wpeframework.service
After=wpeframework.service
ConditionPathExists=/tmp/wpeframeworkstarted
[Service]
Type=oneshot
RemainAfterExit=yes
TimeoutStartSec=30
TimeoutStopSec=10
PrivateTmp=yes
ExecStart=/usr/bin/PluginActivator Monitor

[Install]
WantedBy=multi-user.target
