# systemd/system/wpeframework-messenger.service
# WPEFramework Messenger Initialiser
# Activates the Messenger plugin after framework initialization
[Unit]
Description=WPEFramework Messenger Initialiser
Requires=wpeframework.service
After=wpeframework.service
ConditionPathExists=/tmp/wpeframeworkstarted
[Service]
Type=oneshot
RemainAfterExit=yes
TimeoutStartSec=30
TimeoutStopSec=10
PrivateTmp=yes
ExecStart=/usr/bin/PluginActivator Messenger

[Install]
WantedBy=multi-user.target
