# systemd/system/wpeframework-locationsync.service
# WPEFramework LocationSync Initialiser
# Activates the LocationSync plugin after framework initialization
[Unit]
Description=WPEFramework LocationSync Initialiser
Requires=wpeframework.service
After=wpeframework.service wpeframework-network.service
ConditionPathExists=/tmp/wpeframeworkstarted
[Service]
Type=oneshot
RemainAfterExit=yes
TimeoutStartSec=30
TimeoutStopSec=10
PrivateTmp=yes
ExecStart=/usr/bin/PluginActivator LocationSync

[Install]
WantedBy=multi-user.target
