# systemd/system/wpeframework-launchdelegate.service
[Unit]
Description=WPEFramework LaunchDelegate Initialiser
Requires=wpeframework.service
After=wpeframework.service
ConditionPathExists=/tmp/wpeframeworkstarted

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/PluginActivator org.rdk.LaunchDelegate

