# systemd/system/wpeframework-displaysettings.service
[Unit]
Description=WPEFramework DisplaySettings Initialiser
Requires=dsmgr.service wpeframework-powermanager.service wpeframework-hdmicecsink.service
After=dsmgr.service wpeframework-powermanager.service wpeframework-hdmicecsink.service
[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/PluginActivator org.rdk.DisplaySettings
