# systemd/system/wpeframework-telemetrymetrics.service
[Unit]
Description=WPEFramework TelemetryMetrics Initialiser
Requires= wpeframework.service
After= wpeframework.service
ConditionPathExists=/opt/ai2managers
ConditionPathExists=/tmp/wpeframeworkstarted

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/PluginActivator org.rdk.TelemetryMetrics

