# systemd/system/wpeframework-bluetooth.service
[Unit]
Description=WPEFramework Bluetooth Initialiser
Requires=wpeframework.service iarmbusd.service btmgr.service
After=wpeframework.service iarmbusd.service btmgr.service
[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/PluginActivator org.rdk.Bluetooth
