# systemd/system/wpeframework-usbdevice.service
[Unit]
Description=WPEFramework UsbDevice Initialiser
Requires=wpeframework.service
After=wpeframework.service
ConditionPathExists=/tmp/wpeframeworkstarted
[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/PluginActivator org.rdk.UsbDevice
