# systemd/system/wpeframework-appgateway.service
[Unit]
Description=WPEFramework AppGateway Initialiser
Requires=wpeframework.service wpeframework-launchdelegate.service wpeframework-appnotifications.service wpeframework-app2appprovider.service
After=wpeframework.service wpeframework-launchdelegate.service wpeframework-appnotifications.service wpeframework-app2appprovider.service
ConditionPathExists=/tmp/wpeframeworkstarted
ConditionPathExists=/etc/appgatewayenabled

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/PluginActivator org.rdk.AppGateway
