# systemd/system/wpeframework-provisioning-ready.path
[Unit]
Description=Wait for Thunder PROVISIONING precondition
Documentation=Thunder framework provisioning precondition fulfilled
After=wpeframework-deviceprovisioning.service

[Path]
PathExists=/tmp/thunder-provisioning-precondition
Unit=wpeframework-provisioning-ready.target
MakeDirectory=false

[Install]
WantedBy=multi-user.target