# systemd/system/wpeframework-voicecontrol.service
[Unit]
Description=WPEFramework VoiceControl Initialiser
Requires=wpeframework.service iarmbusd.service ctrlm-main.service
After=wpeframework.service iarmbusd.service ctrlm-main.service
[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/PluginActivator org.rdk.VoiceControl
